# supervisord.conf
[program:node]
command = env PORT=8080 node bin/strider
directory = /home/dotcloud/current
